# autostart/setup.desktop
Type=Application
Name=AutostartScript
Exec=~/.config/autostart/setup.sh
